# .config/systemd/user/reload-waybar.path
[Unit]
Description=Watch waybar for changes

[Path]
PathModified=%h/.config/waybar/style.css
PathModified=%h/.config/waybar/config.jsonc
